# .manala/templates/services/php/run
#!/bin/sh
exec php-fpm --fpm-config /etc/php/default/fpm/project.conf
