# .manala/templates/services/nginx/run
#!/bin/sh
exec nginx -c /etc/nginx/project.conf
